# system_files/etc/udev/rules.d/99-usb-autmount.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", RUN+="/usb-mount.sh %k"
ACTION=="remove", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", RUN+="/usb-unmount.sh %k"